# crontab
# All of your crons will go in here
* * * * * php /var/www/html/artisan schedule:run
